# scripts/nix-scheduler-processor.service
[Unit]
Description=Nix — scheduler output processor (safety net)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=topp
WorkingDirectory=/home/topp/projects/aquarisamatiran
ExecStart=/home/topp/projects/aquarisamatiran/.venv/bin/python -m nixfw process-scheduler-results
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
